# etc/systemd/system/multi-ear-wifi.timer
[Unit]
Description=Multi-EAR Wi-Fi scheduled script execution

[Timer]
OnActiveSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
